# idp-node/conf/supervisord.tomcat.conf
[program:tomcat]
command=sh -c "/usr/local/bin/wait_for_postgres.sh && %(ENV_CATALINA_HOME)s/bin/catalina.sh run"
process_name=%(program_name)s
startsecs=5
stopsignal=INT
user=tomcat
#stdout_logfile=%(ENV_CATALINA_HOME)s/logs/catalina.out
#stderr_logfile=%(ENV_CATALINA_HOME)s/logs/catalina.err
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
