# postgres/conf/supervisord.postgres.conf
; Supervisor configuration for Postgres database.

[program:postgres]
command=/bin/bash -c "/usr/bin/postgres -D /var/lib/pgsql/data -p 5432"
user=postgres
