# dashboard/conf/supervisord.dashboard.conf
[program:dashboard]

command=%(ENV_DASHBOARD_HOME)s/bin/esgf-dashboard-ip
process_name=%(program_name)s
startsecs=5
user=root
numprocs=1
stdout_logfile=/tmp/esgf_dashboard.log
redirect_stderr=true
