# solr/conf/supervisord.solr_master_8984.conf
[program:master_8984]
user=solr
command=%(ENV_SOLR_INSTALL_DIR)s/bin/solr start -f -d %(ENV_SOLR_INSTALL_DIR)s/server -s %(ENV_SOLR_HOME)s/master-8984 -p 8984 -a '-Denable.master=true'
