# solr/conf/supervisord.solr_slave_8983.conf
[program:slave_8983]
user=solr
command=%(ENV_SOLR_INSTALL_DIR)s/bin/solr start -f -d %(ENV_SOLR_INSTALL_DIR)s/server -s %(ENV_SOLR_HOME)s/slave-8983 -p 8983 -a '-Denable.slave=true -Denable.master=true'
