# data-node/dashboard/conf/supervisord.dashboard.conf
[program:dashboard]
command=%(ENV_DASHBOARD_HOME)s/bin/ip.service start
process_name=%(program_name)s
startsecs=5
user=dashboard
numprocs=1 
stdout_logfile=%(ENV_DASHBOARD_HOME)s/logs/esgf_dashboard.log
redirect_stderr=true
